# client/linux/packaging/debian/prerm
#!/bin/sh
set -e

ACTION="${1:-}"

stop_for_user() {
  user="$1"
  uid="$2"
  runtime="/run/user/$uid"

  if [ ! -d "$runtime" ]; then
    return 0
  fi

  su - "$user" -s /bin/sh -c "XDG_RUNTIME_DIR=$runtime systemctl --user disable --now virtue.service >/dev/null 2>&1 || true"
  su - "$user" -s /bin/sh -c "XDG_RUNTIME_DIR=$runtime systemctl --user disable --now bepure.service >/dev/null 2>&1 || true"
}

# Keep global enable on upgrades; remove it on remove/deconfigure.
if [ "$ACTION" = "remove" ] || [ "$ACTION" = "deconfigure" ]; then
  systemctl --global disable virtue.service >/dev/null 2>&1 || true
fi

if command -v loginctl >/dev/null 2>&1; then
  loginctl list-users --no-legend | while read -r uid user _; do
    [ -n "$uid" ] || continue
    [ -n "$user" ] || continue
    stop_for_user "$user" "$uid"
  done
fi

if [ -n "$SUDO_USER" ] && [ "$SUDO_USER" != "root" ]; then
  uid=$(id -u "$SUDO_USER" 2>/dev/null || true)
  if [ -n "$uid" ]; then
    stop_for_user "$SUDO_USER" "$uid"
  fi
fi

exit 0
